# ngsi_adapter/tools/build/files/debian/postrm
#!/bin/sh
# postrm script for fiware-monitoring-ngsi-adapter

set -e

# actual values of these variables will be replaced after dh_install (see rules)
ADAPTER_SRV=#DH_ADAPTER_SRV#
ADAPTER_DIR=#DH_ADAPTER_DIR#
LOGGING_DIR=#DH_LOGGING_DIR#
FIWARE_DIR=#DH_FIWARE_DIR#

case "$1" in
remove)
	# remove Linux service
	update-rc.d -f $ADAPTER_SRV remove

	# remove installation directory
	rm -rf $ADAPTER_DIR

	# remove FIWARE parent directory, if empty
	[ -d $FIWARE_DIR ] && rmdir --ignore-fail-on-non-empty $FIWARE_DIR
	;;

purge)
	# remove log files
	rm -rf $LOGGING_DIR
	;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
